# images/arm64/nxp/frdm_imx93/systemd/system_config/etc/systemd/system/tee-supplicant.service
[Unit]
Description=TEE Supplicant
After=optee.service
Wants=optee.service

[Service]
Type=simple
ExecStart=/usr/sbin/tee-supplicant
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
