# images/qemu-crinit-x86_64/root/etc/profile.d/10-ash-prompt.sh
if [ ${SHELL} == "/bin/ash" ]; then
    PS1='\u:\w\$ '
fi
