# images/nxp-s32g2-systemd/root/etc/systemd/system/data.mount
[Unit]
Description=Mount Data Partition
Before=multi-user.target
[Mount]
What=/dev/vda2
Where=/data
Type=ext4
Options=rw
