# images/rdb2-systemd/root/etc/systemd/system/serial-getty@ttyLF0.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/sbin/agetty --noclear --keep-baud 115200,38400,9600 ttyLF0 $TERM
