# images/amd64/qemu/xen/dom0/systemd/systemd_config/etc/systemd/system/start-ebcl-1.service
[Unit]
Description=Start Ebcl-1 VM
After=multi-user.target

[Service]
ExecStart=/usr/sbin/xl create -f /etc/xen/ebcl.cfg
Type=simple

[Install]
WantedBy=multi-user.target
